# rpicam-record-proccess.conf
[program:rpicam-record-proccess]
directory=/home/pi/camera_record
#command=python3 ./camera_server.py
command=python3 ./record.py
autostart=true
autorestart=true
user=pi

